# deploy/huas-api.service
[Unit]
Description=HUAS API Service
Documentation=https://github.com/your-repo/huas-api
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/huas-api
Environment="NODE_ENV=production"
Environment="PORT=12103"
# 监控端口（默认 13001），如需关闭设为 0
Environment="MONITOR_PORT=13001"
# 监控监听地址（默认 0.0.0.0，若仅本机访问可设 127.0.0.1）
Environment="MONITOR_HOST=0.0.0.0"
EnvironmentFile=/opt/huas-api/.env.production
ExecStart=/usr/local/bin/bun run start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huas-api

# 安全加固
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/huas-api /var/log/huas-api

# 允许以非 root 绑定 80/443 等低端口（需要对 bun 设置 setcap）
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# 资源限制
LimitNOFILE=65536
LimitNPROC=512

[Install]
WantedBy=multi-user.target
